# udev/rules.d/70-snap.canonical-livepatch.rules
# This file is automatically generated.
# kernel-module-control
KERNEL=="mem", TAG+="snap_canonical-livepatch_canonical-livepatch"
# network-control
KERNEL=="rfkill", TAG+="snap_canonical-livepatch_canonical-livepatch"
# network-control
KERNEL=="tun", TAG+="snap_canonical-livepatch_canonical-livepatch"
TAG=="snap_canonical-livepatch_canonical-livepatch", RUN+="/usr/lib/snapd/snap-device-helper $env{ACTION} snap_canonical-livepatch_canonical-livepatch $devpath $major:$minor"
# kernel-module-control
KERNEL=="mem", TAG+="snap_canonical-livepatch_canonical-livepatchd"
# network-control
KERNEL=="rfkill", TAG+="snap_canonical-livepatch_canonical-livepatchd"
# network-control
KERNEL=="tun", TAG+="snap_canonical-livepatch_canonical-livepatchd"
TAG=="snap_canonical-livepatch_canonical-livepatchd", RUN+="/usr/lib/snapd/snap-device-helper $env{ACTION} snap_canonical-livepatch_canonical-livepatchd $devpath $major:$minor"
